# docker/chrome/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[program:Xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:fluxbox]
command=fluxbox -display :99
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
depends_on=Xvfb

[program:x11vnc]
command=x11vnc -display :99 -rfbport 5900 -rfbauth /root/.vnc/passwd -forever -shared -nopw
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
depends_on=Xvfb

[program:chrome]
command=google-chrome-stable
    --no-sandbox
    --disable-gpu
    --disable-dev-shm-usage
    --remote-debugging-port=9222
    --remote-debugging-address=0.0.0.0
    --window-size=1920,1080
    --display=:99
    --user-data-dir=/tmp/chrome-data
    --disable-background-timer-throttling
    --disable-backgrounding-occluded-windows
    --disable-renderer-backgrounding
    about:blank
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40
depends_on=Xvfb
